# build/etc/systemd/system/scraper.service
[Unit]
Description=repeater sensors scraper

[Service]
ExecStart=/opt/repeater-sensors/scraper
WorkingDirectory=/opt/repeater-sensors
Restart=always

[Install]
WantedBy=multi-user.target
